# power-monitor/.config/systemd/user/power-monitor.service
[Service]
ExecStart=/home/lspita/.local/bin/power-monitor.sh

[Install]
WantedBy=default.target